# home/dot_config/systemd/user/dstask-notify.timer.tmpl
{{ if and (lookPath "dstask") (lookPath .notifier) -}}
[Unit]
Description=dstask due-today digest — weekdays at 09:30

[Timer]
OnCalendar=Mon-Fri *-*-* 09:30:00
Persistent=true
AccuracySec=1min

[Install]
WantedBy=timers.target
{{- end }}
