# home/dot_config/systemd/user/rano-weekly-report.service.tmpl
{{ if lookPath "rano" -}}
[Unit]
Description=rano weekly activity report (disk, CPU, provider attribution)
Documentation=file://%h/.local/share/chezmoi/home/dot_local/bin/executable_rano-weekly-report.tmpl
Requires=rano-monitor.service
After=rano-monitor.service

[Service]
Type=oneshot
ExecStart=%h{{ trimPrefix .chezmoi.homeDir (output "chezmoi" "target-path" (joinPath .chezmoi.sourceDir "dot_local/bin/executable_rano-weekly-report.tmpl") | trim) }}
{{- end }}
